# services/nagios/run
#!/bin/sh

exec /opt/nagios/bin/nagios /opt/nagios/etc/nagios.cfg > /var/log/nagios.log

